# systemd/hermes-memory-stack.system.target
[Unit]
Description=Hermes Memory Stack (PostgreSQL + Hindsight; Memory Graph may be user-unit owned)
Documentation=https://hermes-agent.nousresearch.com/docs
Wants=postgresql@15-main.service hindsight.service
After=postgresql@15-main.service hindsight.service

[Install]
WantedBy=multi-user.target
