# systemd/hermes-memory-stack.system.target
[Unit]
Description=Hermes Memory Stack (PostgreSQL + Hindsight + Memory Graph)
Documentation=https://hermes-agent.nousresearch.com/docs
Wants=postgresql@15-main.service hindsight.service hermes-memory-graph.service
After=postgresql@15-main.service hindsight.service hermes-memory-graph.service

[Install]
WantedBy=multi-user.target
